# wallfacer.service
[Unit]
Description=Wallfacer demoscene effects and projection mapping
After=multi-user.target

[Service]
Type=simple
User=root
WorkingDirectory=/home/mapper/wallfacer
Environment=SDL_VIDEODRIVER=kmsdrm
Environment=SDL_KMSDRM_DEVICE=/dev/dri/card1
Environment=HOME=/home/mapper
# Stop getty on tty1 to release display before starting
ExecStartPre=/bin/systemctl stop getty@tty1.service
ExecStart=/home/mapper/wallfacer/target/release/wallfacer -r 640x480
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
