# wallfacer.service
[Unit]
Description=Wallfacer demoscene effects and projection mapping
After=multi-user.target
Conflicts=getty@tty1.service

[Service]
Type=simple
User=root
WorkingDirectory=/home/mapper/wallfacer
Environment=SDL_VIDEODRIVER=kmsdrm
Environment=SDL_KMSDRM_DEVICE=/dev/dri/card1
Environment=HOME=/home/mapper
ExecStart=/home/mapper/wallfacer/target/release/wallfacer -r 1920x1080
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
